# config/crontab
* * * * * cd /var/www/html && php artisan schedule:run >> /dev/null 2>&1
# cron needs an empty line at the end of the file to work
